# supervisord.conf
[supervisord]
nodaemon=true

[program:streama]
command=/opt/streama.war
redirect_stderr=true
stdout_logfile=/var/log/streama.log
stopsignal=TERM
autostart=true
autorestart=true